# deploy/systemd/codex-agent-backend.service
[Unit]
Description=Codex Agent Backend (API-Only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/codex_agent
EnvironmentFile=%h/codex_agent/deploy/codex-backend.env
ExecStartPre=/usr/bin/mkdir -p %h/codex_agent/workspace/.agent_state/logs
ExecStart=%h/codex_agent/deploy/scripts/run_codex_backend.sh
Restart=always
RestartSec=3
StandardOutput=append:%h/codex_agent/workspace/.agent_state/logs/codex-backend.log
StandardError=append:%h/codex_agent/workspace/.agent_state/logs/codex-backend.error.log

[Install]
WantedBy=default.target
